# meta-riscv/recipes-core/mount-overlayfs/mount-overlayfs/mount-overlayfs.service
[Unit]
Description=Overylay mount
#DefaultDependencies=no
After=systemd-networkd.service

[Service]
Type=simple
ExecStart=/etc/nezha/mount_root.sh

[Install]
WantedBy=multi-user.target
