# meta-riscv/recipes-bsp/d1/d1-nezha-init/d1-nezha-init.service
[Unit]
Description=D1 Init Service

[Service]
Type=simple
ExecStart=/bin/sh /etc/nezha/d1-nezha-init.sh

[Install]
WantedBy=multi-user.target
